# docs/com.kbo.title-news-automation.watcher.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kbo.title-news-automation.watcher</string>

  <key>WorkingDirectory</key>
  <string>/Users/s.h.choi/Desktop/kbo/title_automation/title news code</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>PYTHONDONTWRITEBYTECODE=1</string>
    <string>PYTHONPATH=src</string>
    <string>python</string>
    <string>-m</string>
    <string>kbo_card_news.automation.cli</string>
    <string>watch-cycle</string>
    <string>--candidate-count</string>
    <string>10</string>
    <string>--max-candidates</string>
    <string>10</string>
    <string>--notify</string>
  </array>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>StandardOutPath</key>
  <string>/Users/s.h.choi/Desktop/kbo/title_automation/title news code/outputs/automation/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/s.h.choi/Desktop/kbo/title_automation/title news code/outputs/automation/logs/launchd.err.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
